# files/etc/crontabs/root
# weekly reboot
12 3 * * 3 sleep 70 && touch /etc/banner && reboot
